# valorant-bot.service
# ─────────────────────────────────────────────────────────────────────────────
#  valorant-bot.service  —  systemd unit for the VALORANT Tournament Discord Bot
#
#  INSTALLATION (on Raspberry Pi):
#  1. Copy this file to /etc/systemd/system/valorant-bot.service
#  2. Edit the paths below to match your actual setup
#  3. sudo systemctl daemon-reload
#  4. sudo systemctl enable valorant-bot   ← start on boot
#  5. sudo systemctl start  valorant-bot   ← start right now
#
#  USEFUL COMMANDS:
#    sudo systemctl status  valorant-bot   ← check if it's running
#    sudo systemctl restart valorant-bot   ← restart after code changes
#    sudo systemctl stop    valorant-bot   ← stop the bot
#    journalctl -u valorant-bot -f         ← live logs
# ─────────────────────────────────────────────────────────────────────────────

[Unit]
Description=VALORANT PC Tournament Discord Bot
Documentation=https://discordpy.readthedocs.io
After=network-online.target
Wants=network-online.target
# Restart the unit if the network goes down and comes back
StartLimitIntervalSec=120
StartLimitBurst=5

[Service]
# ── Change these two paths to match your Pi setup ────────────────────────────
# The user that will own the process (avoid running as root!)
User=pi
Group=pi

# Absolute path to your project directory
WorkingDirectory=/home/pi/valorant-bot

# Absolute path to the Python interpreter inside your venv
ExecStart=/home/pi/valorant-bot/venv/bin/python bot.py
# ─────────────────────────────────────────────────────────────────────────────

# Gracefully stop: send SIGTERM, wait 15 s, then SIGKILL
KillSignal=SIGTERM
TimeoutStopSec=15

# Always restart on failure (crashes, OOM, etc.)
Restart=always
RestartSec=10

# Environment hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/pi/valorant-bot   # the bot writes bot.log here

# Stdout/stderr go to the journal (view with: journalctl -u valorant-bot -f)
StandardOutput=journal
StandardError=journal
SyslogIdentifier=valorant-bot

[Install]
WantedBy=multi-user.target
